# phpService/mytemp.service
[Unit]
Description=Demo php service

[Service]
ExecStart=/home/ps/phpService/mytemp.php

[Install]
WantedBy=multi-user.target